# snapserver/rootfs/etc/services.d/vlc/run

#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start VLC service
# ==============================================================================

# Run daemon
exec cvlc -I telnet --telnet-port 4212 --telnet-password testest --no-video --aout afile --audiofile-file /share/snapfifo/vlc-tous
